# tools/sysd-service/roverforge.service
[Unit]
Description=RoverForge Pi Control Server
After=network.target

[Service]
Type=simple
User=larry
WorkingDirectory=/home/larry/Development/RoverForge/pi_control
# Use the virtualenv python if available; venv is created by the developer start script.
ExecStart=/home/larry/.venvs/roverforge/bin/python /home/larry/Development/RoverForge/pi_control/server.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Optional: set default Arduino serial and baud via environment
Environment=ARDUINO_PORT=/dev/ttyACM0
Environment=ARDUINO_BAUD=9600

[Install]
WantedBy=multi-user.target
